# opengateway-ai-proxy.service
[Unit]
Description=AI Proxy (OpenAI-compatible Flask gateway)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ai-proxy
EnvironmentFile=/opt/ai-proxy/.env
ExecStart=/opt/ai-proxy/venv/bin/gunicorn -c /opt/ai-proxy/gunicorn.conf.py "app:create_app()"
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

# Allow binding to port 80 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging via journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-proxy

[Install]
WantedBy=multi-user.target